# conf.d/soapui.conf
[program:SoapUI]
command=SoapUI-5.7.0
autorestart=false
